# etc/udev/rules.d/10-amdgpu.rules
ACTION!="remove", SUBSYSTEM=="drm", KERNEL=="card0", DRIVERS=="amdgpu", RUN+="/usr/local/bin/scripts/amdgpufix.sh"
